# etc-gentoo/profile.d/my.sh
#!/bin/sh

export LANG=en_US.UTF-8

alias df='df -h'
alias f='free -m'
alias ll='ls -la'
alias m='mount '

alias es='eix-sync && eix-update'
alias eu='etc-update && env-update && . /etc/profile'
alias ew='emerge -uDNpv world'

